# com.francoisdupont.briefing-assistant.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.francoisdupont.briefing-assistant</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/docker</string>
    <string>compose</string>
    <string>--project-directory</string>
    <string>/Users/francoisdupont/.briefing-assistant</string>
    <string>up</string>
    <string>--remove-orphans</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <dict>
    <key>SuccessfulExit</key>
    <false/>
  </dict>

  <key>StartInterval</key>
  <integer>10</integer>

  <key>StandardOutPath</key>
  <string>/Users/francoisdupont/.briefing-assistant/logs/briefing.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/francoisdupont/.briefing-assistant/logs/briefing.err</string>

  <key>ThrottleInterval</key>
  <integer>30</integer>
</dict>
</plist>
